# jenkins-setup/jenkins-ec2-setup/ansible/roles/jenkins/templates/jenkins.service.j2
[Unit]
Description=Jenkins Continuous Integration Server
After=network.target

[Service]
User=jenkins
Environment=JAVA_HOME=/usr/lib/jvm/java-11-openjdk-amd64
ExecStart=/usr/bin/java -jar /usr/share/jenkins/jenkins.war
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target